# prisonguard.target
[Unit]
Description=Prison Guard Online
Requires=prison-ns.service guard-ns.service prisonguard-interface.service
After=prisonguard-interface.service
